# config/supervisor/uwsgi.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
